# infra/gcp/stage1/lumberjacks-compose.service
[Unit]
Description=Lumberjacks Gate 1 Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists=/opt/lumberjacks/infra/docker/docker-compose.yml

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/lumberjacks
EnvironmentFile=/etc/lumberjacks/environment
ExecStart=/usr/bin/docker compose -f infra/docker/docker-compose.yml -f infra/docker/docker-compose.gcp-stage1.yml up -d
ExecStop=/usr/bin/docker compose -f infra/docker/docker-compose.yml -f infra/docker/docker-compose.gcp-stage1.yml down
TimeoutStartSec=900
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
